# shards/test.service
[Unit]
Description=Test Service

[Service]
Type=forking
ExecStart=/bin/bash -c 'sleep 10'

[Install]
WantedBy=multi-user.target
